# projects/nexus/deploy/systemd/mausoleum.service
[Unit]
Description=Mausoleum Document Database
Documentation=https://ritz-lang.org/mausoleum
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ritz
Group=ritz

# Binary location
ExecStart=/opt/ritz/bin/mausoleum serve \
    --port 7777 \
    --data-dir /var/lib/mausoleum \
    --no-encryption

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/mausoleum /var/log/ritz

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mausoleum

# Resource limits
LimitNOFILE=65536
MemoryMax=256M

[Install]
WantedBy=multi-user.target
