# projects/nexus/deploy/systemd/valet.service
[Unit]
Description=Valet HTTPS Server (Frontend for Zeus)
Documentation=https://ritz-lang.org/valet
After=network.target zeus.service
Requires=zeus.service

[Service]
Type=simple
User=root
Group=root

# Valet HTTPS server with Zeus backend
# -c: Configuration file with TLS settings
# -p 443: Listen on HTTPS port (overrides config)
# --zeus-socket: Connect to Zeus daemon for worker dispatch
ExecStart=/opt/ritz/bin/valet -c /etc/ritz/valet.json -p 443 --zeus-socket /run/zeus/zeus.sock

# Restart policy
Restart=on-failure
RestartSec=5

# Security (limited due to privileged port requirement)
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=valet

# Resource limits - needs higher for io_uring
LimitNOFILE=65536
LimitMEMLOCK=infinity
MemoryMax=256M

# Capabilities for binding to ports 80/443
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
