# projects/nexus/deploy/systemd/valet.service
[Unit]
Description=Valet HTTP Server (Frontend for Zeus)
Documentation=https://ritz-lang.org/valet
After=network.target zeus.service
Requires=zeus.service

[Service]
Type=simple
User=root
Group=root

# Valet HTTP server
# -p 80: Listen on port 80
# -w 1: Single process (zeus handles workers)
# --zeus-socket: Connect to Zeus daemon for worker dispatch
ExecStart=/opt/ritz/bin/valet -p 80 -w 1 --zeus-socket /run/zeus/zeus.sock

# Restart policy
Restart=on-failure
RestartSec=5

# Security (limited due to port 80 requirement)
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=valet

# Resource limits - needs higher for io_uring
LimitNOFILE=65536
LimitMEMLOCK=infinity
MemoryMax=256M

# Capabilities for binding to port 80
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
