# projects/nexus/deploy/systemd/zeus.service
[Unit]
Description=Zeus Daemon (Nexus Workers)
Documentation=https://ritz-lang.org/zeus
After=network.target mausoleum.service
Requires=mausoleum.service

[Service]
Type=simple
User=root
Group=root

# Environment for mausoleum connection (passed to workers)
Environment=MAUSOLEUM_HOST=127.0.0.1
Environment=MAUSOLEUM_PORT=7777

# Create runtime directory for zeus socket
RuntimeDirectory=zeus
RuntimeDirectoryMode=0755

# Zeus daemon with nexus worker
# -w 2: 2 worker processes
# -m 2097152: 2MB shared memory per worker
# -f: foreground (don't daemonize, let systemd manage)
# -s: Unix socket path for Valet connections
ExecStart=/opt/ritz/bin/zeus -w 2 -m 2097152 -f -s /run/zeus/zeus.sock /opt/ritz/bin/nexus

# Restart policy
Restart=on-failure
RestartSec=5

# Security
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zeus

# Resource limits
LimitNOFILE=65536
LimitMEMLOCK=infinity
MemoryMax=512M

[Install]
WantedBy=multi-user.target
